# docker/supervisord.conf
[supervisord]
nodaemon=true
user=root

[program:postgresql]
command=/usr/lib/postgresql/14/bin/postgres -D /var/lib/postgresql/14/main -c config_file=/etc/postgresql/14/main/postgresql.conf
user=postgres
autostart=true
autorestart=true
priority=1
stdout_logfile=/var/log/supervisor/postgresql.log
stderr_logfile=/var/log/supervisor/postgresql.log

[program:redis]
command=redis-server /etc/redis/redis.conf
autostart=true
autorestart=true
priority=2
stdout_logfile=/var/log/supervisor/redis.log
stderr_logfile=/var/log/supervisor/redis.log

[program:init-db]
command=/app/init-db.sh
autostart=true
autorestart=false
priority=3
startsecs=0
exitcodes=0
stdout_logfile=/var/log/supervisor/init-db.log
stderr_logfile=/var/log/supervisor/init-db.log

[program:secman-app]
command=/app/start-backend.sh
directory=/app
autostart=true
autorestart=true
priority=4
environment=DB_HOST=localhost,DB_PORT=5432,DB_USER=lutenix,DB_PASSWORD=lutenix,DB_NAME=lutenix_db,REDIS_URL=localhost:6379
stdout_logfile=/var/log/supervisor/secman-app.log
stderr_logfile=/var/log/supervisor/secman-app.log

[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
autostart=true
autorestart=true
priority=5
stdout_logfile=/var/log/supervisor/nginx.log
stderr_logfile=/var/log/supervisor/nginx.log